# debian/src/deb/control/prerm
#!/usr/bin/env bash

for filename in /usr/share/corfu/bin/*; do
    rm /usr/local/bin/$(basename "$filename")
done

update-rc.d -f corfu-server remove
