# debian/src/deb/control-inmem/prerm
#!/usr/bin/env bash

if [[ -n $(which systemctl) ]]; then
    systemctl disable corfu-server-inmem
else
    update-rc.d corfu-server-inmem remove
fi
